# meta-agl-ic-container/recipes-demo/momikey/files/momikey.rules
SUBSYSTEM=="usb", ATTR{idVendor}=="1189", ATTR{idProduct}=="8890", RUN+="/bin/sh -c 'systemctl start momikey'"
